# vm/etc/profile.d/tango.sh
# adding TANGO_HOST
. /etc/tangorc
export TANGO_HOST
export TANGO_ROOT
export ARCHIVING_ROOT
export HDBPP_ROOT
