# docker/base/resources/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/tmp/supervisord.log
pidfile=/var/tmp/supervisord.pid

[program:tango-starter]
command=/bin/bash -c '
  /usr/local/bin/tango_register_device Starter/`hostname --short` Starter tango/admin/`hostname --short` &&
  exec -- /usr/bin/Starter `hostname --short`'
autorestart=true
startretries=4
priority=1

[program:tango-base]
command=/bin/bash -c '
  exec -- python /usr/local/bin/PyAlarm.py Example_Alarm_Base -v4'
autorestart=true
startretries=5
priority=2
