# docker/base/resources/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/tmp/supervisord.log
pidfile=/var/tmp/supervisord.pid

[program:tango-starter]
command=/bin/bash -c '
  chmod 755 /usr/local/bin/tango_register_device &&
  /usr/local/bin/tango_register_device Starter/`hostname --short` Starter tango/admin/`hostname --short` &&
  exec -- /usr/local/bin/Starter `hostname --short`'
autorestart=true
startsecs=0
startretries=5
priority=1

[program:tango-base]
command=/bin/bash -c '
  chmod 755 /usr/local/lib/python2.7/dist-packages/panic/ds/PyAlarm.py &&
  exec -- python /usr/local/lib/python2.7/dist-packages/panic/ds/PyAlarm.py Example_Alarm_`hostname --short` -v4 -D'
autorestart=true
startsecs=0
priority=2
